# scriptsKDE/ChangeWallpaper/changewallpaper.timer
[Unit]
Description=Change wallpaper every minute

[Timer]
OnUnitActiveSec=1min
OnBootSec=1min

[Install]
WantedBy=timers.target